# conffiles/config/systemd/user/awesome.service
[Unit]
Description=Awesome window manager

[Service]
WorkingDirectory=/home/thequux
ExecStart=/usr/bin/awesome

[Install]
WantedBy=graphical.target
