# conffiles/config/systemd/user/chromium.service
[Unit]
Description=Chromium

[Service]
ExecStart=/usr/bin/chromium

[Install]
WantedBy=graphical.target
