# conffiles/config/systemd/user/gmpc.service
[Unit]
Description=Gnome media player client

[Service]
ExecStart=/usr/bin/gmpc

[Install]
WantedBy=graphical.target
