# conffiles/config/systemd/user/graphical.target
[Unit]
Description=Graphical desktop
Wants=common.target
