# conffiles/config/systemd/user/iceweasel.service
[Unit]
Description=Iceweasel

[Service]
ExecStart=/usr/local/bin/firefox

[Install]
WantedBy=graphical.target
