# conffiles/config/systemd/user/irssi.service
[Unit]
Description=IRSSI IRC client

[Service]
ExecStart=/usr/bin/x-terminal-emulator -e irssi

[Install]
WantedBy=graphical.target
