# conffiles/config/systemd/user/kontact.service
[Unit]
Description=KDE PIM

[Service]
ExecStart=/usr/bin/kontact --nofork

[Install]
WantedBy=graphical.target
