# conffiles/config/systemd/user/nepomuk.service
[Unit]
Description=KDE Semantic Desktop

[Service]
Type=forking
ExecStart=/usr/bin/nepomukctl start

[Install]
WantedBy=graphical.target
