# conffiles/config/systemd/user/nm-applet.service
[Unit]
Description=Network manager applet

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=graphical.target
