# conffiles/config/systemd/user/pidgin.service
[Unit]
Description=Pidgin

[Service]
ExecStart=/usr/bin/pidgin

[Install]
WantedBy=graphical.target
