# conffiles/config/systemd/user/psi.service
[Unit]
Description=Psi

[Service]
ExecStart=/usr/bin/psi-plus

[Install]
WantedBy=graphical.target
