# conffiles/config/systemd/user/template.service
[Unit]
Description=

[Service]
ExecStart=/usr/bin/emacs --daemon

[Install]
WantedBy=graphical.target
