# conffiles/config/systemd/user/wicd-gtk.service
[Unit]
Description=Network manager applet

[Service]
ExecStart=/usr/bin/wicd-gtk -t

[Install]
WantedBy=graphical.target
