# conffiles/config/systemd/user/xbattbar.service
[Unit]
Description=Xbattbar

[Service]
ExecStart=/usr/bin/xbattbar -c -a

[Install]
WantedBy=graphical.target
