# conffiles/config/systemd/user/xchat.service
[Unit]
Description=XChat

[Service]
ExecStart=/usr/bin/xchat

[Install]
WantedBy=graphical.target
