# conffiles/config/systemd/user/xrdb.service
[Unit]
Description=Xrdb

[Service]
Type=oneshot
ExecStart=/usr/bin/xrdb /home/thequux/.Xdefaults

[Install]
WantedBy=graphical.target
